# roles/editor/templates/galaxy.service.j2
# /lib/systemd/system/galaxy.service
[Unit]
Description=Galaxy Workflow Editor

[Service]
User=root
Type=simple
Restart=on-failure
ExecStart={{ editor_directory }}/run.sh